# deploy/supervisor.conf
; Meeting Room Booking System - Supervisor Configuration
; Production-ready configuration for background workers
;
; @author Digital Technology & AI Division
; @version 1.0.0
;
; Installation:
; 1. Copy to /etc/supervisor/conf.d/meeting-room-booking.conf
; 2. Run: supervisorctl reread
; 3. Run: supervisorctl update
; 4. Run: supervisorctl start all

; ========================================
; Email Queue Worker
; Processes email sending queue
; ========================================
[program:booking-email-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/meeting-room-booking/yii queue/listen --verbose=1 --color=0
directory=/var/www/meeting-room-booking
autostart=true
autorestart=true
startretries=10
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/supervisor/booking-email-worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stopwaitsecs=30
stopsignal=SIGTERM

; ========================================
; Notification Worker
; Processes notification queue
; ========================================
[program:booking-notification-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/meeting-room-booking/yii notification-queue/listen --verbose=1 --color=0
directory=/var/www/meeting-room-booking
autostart=true
autorestart=true
startretries=10
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/booking-notification-worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stopwaitsecs=30
stopsignal=SIGTERM

; ========================================
; Scheduler (Cron replacement)
; Runs scheduled tasks every minute
; ========================================
[program:booking-scheduler]
command=php /var/www/meeting-room-booking/yii schedule/run
directory=/var/www/meeting-room-booking
autostart=true
autorestart=true
startretries=3
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/booking-scheduler.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3
stopwaitsecs=10

; ========================================
; Daily Report Generator
; Generates and sends daily reports
; ========================================
[program:booking-daily-report]
command=php /var/www/meeting-room-booking/yii daily-report/generate
directory=/var/www/meeting-room-booking
autostart=false
autorestart=unexpected
startretries=3
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/booking-daily-report.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3
exitcodes=0

; ========================================
; Cleanup Worker
; Runs cleanup tasks periodically
; ========================================
[program:booking-cleanup]
command=php /var/www/meeting-room-booking/yii cleanup/run
directory=/var/www/meeting-room-booking
autostart=false
autorestart=unexpected
startretries=3
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/booking-cleanup.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3
exitcodes=0

; ========================================
; Reminder Sender
; Sends booking reminders
; ========================================
[program:booking-reminders]
command=php /var/www/meeting-room-booking/yii send-reminders/run
directory=/var/www/meeting-room-booking
autostart=true
autorestart=true
startretries=3
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/booking-reminders.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3
stopwaitsecs=10

; ========================================
; Auto Complete Bookings
; Marks past bookings as completed
; ========================================
[program:booking-auto-complete]
command=php /var/www/meeting-room-booking/yii auto-complete/run
directory=/var/www/meeting-room-booking
autostart=true
autorestart=true
startretries=3
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/booking-auto-complete.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3
stopwaitsecs=10

; ========================================
; Auto Cancel Bookings
; Cancels stale pending bookings
; ========================================
[program:booking-auto-cancel]
command=php /var/www/meeting-room-booking/yii auto-cancel/run
directory=/var/www/meeting-room-booking
autostart=true
autorestart=true
startretries=3
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/booking-auto-cancel.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3
stopwaitsecs=10

; ========================================
; WebSocket Server (Optional)
; Real-time notifications via WebSocket
; ========================================
[program:booking-websocket]
command=php /var/www/meeting-room-booking/yii websocket/start
directory=/var/www/meeting-room-booking
autostart=false
autorestart=true
startretries=10
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/booking-websocket.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stopwaitsecs=30
stopsignal=SIGTERM

; ========================================
; Process Groups
; ========================================
[group:booking-workers]
programs=booking-email-worker,booking-notification-worker
priority=999

[group:booking-schedulers]
programs=booking-scheduler,booking-reminders,booking-auto-complete,booking-auto-cancel
priority=998

[group:booking-optional]
programs=booking-websocket,booking-daily-report,booking-cleanup
priority=997
